# lights.service
[Unit]
Description=lights service

[Service]
ExecStart=/home/pi/lights/venv/bin/python /home/pi/lights/Client.py

[Install]
WantedBy=multi-user.target
