# lights.service
[Unit]
Description=lights service

[Service]
ExecStart=sh /home/pi/lights/start.sh

[Install]
WantedBy=multi-user.target
